# meshtastic-router.service
[Unit]
Description=Meshtastic Router Service
After=network.target mosquitto.service meshtastic-collector.service
Wants=mosquitto.service

[Service]
Type=simple
User=meshtastic
Group=meshtastic
WorkingDirectory=/opt/meshtastic
Environment=PYTHONPATH=/opt/meshtastic
ExecStart=/usr/bin/python3 -m services.recommendation_engine.router
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target